# systemd/apex-sync-service.service
[Unit]
Description=Apex Sync Service
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=notify
Restart=on-failure
RestartSec=10s
WatchdogSec=60s
# Assume container is named 'sync-service'
ExecStart=/usr/bin/docker start -a sync-service
ExecStop=/usr/bin/docker stop sync-service

# Logging
StandardOutput=journal
StandardError=journal

# Notifications
ExecStartPost=/home/samson/apex/scripts/systemd/telegram_notify.sh "🟢 Sync Service online"
ExecStopPost=/home/samson/apex/scripts/systemd/telegram_notify.sh "⚠️ Sync Service status: $SERVICE_RESULT (Exit Code: $EXIT_STATUS)"

[Install]
WantedBy=multi-user.target
