# systemd/apex-security-warden.service
[Unit]
Description=Apex Security Warden Bot
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=notify
Restart=on-failure
RestartSec=10s
WatchdogSec=60s
# Assume container is named 'security-warden'
ExecStart=/usr/bin/docker start -a security-warden
ExecStop=/usr/bin/docker stop security-warden

# Logging
StandardOutput=journal
StandardError=journal

# Notifications
ExecStartPost=/home/samson/apex/scripts/systemd/telegram_notify.sh "🟢 Security Warden Bot online"
ExecStopPost=/home/samson/apex/scripts/systemd/telegram_notify.sh "⚠️ Security Warden Bot status: $SERVICE_RESULT (Exit Code: $EXIT_STATUS)"

[Install]
WantedBy=multi-user.target
